# cron/crontab
# START CRON JOB
*/1 * * * * scraperuser /home/scraperuser/app/run_scraper.sh
# END CRON JOB